# docs/debian/systemd/eebus-go-controlbox-app.service
# eebus-go-controlbox-app.service
#

[Unit]
Description=eebus-go-controlbox-app
Requires=network-online.target
After=syslog.target network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=10
StartLimitBurst=10

[Service]
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/npm run dev --prefix /opt/eebus-go-apps/apps/controlbox
#Environment="PATH=$PATH:/usr/local/go/bin"
Restart=always
RestartSec=10

User=root
Group=root

[Install]
WantedBy=multi-user.target